# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
stdout_logfile=/var/log/app.log
command=python3 /app/app.py


[program:api]
directory=/app
stdout_logfile=/var/log/api.log
command=python3 /app/api.py